# deploy/systemd/qdgz300-spool-mover.service
[Unit]
Description=QDGZ300 PCAP spool archive mover
After=network-online.target qdgz300-receiver.service
Wants=network-online.target qdgz300-receiver.service

[Service]
Type=simple
User=qdgz300
Group=qdgz300
WorkingDirectory=/opt/qdgz300_backend
ExecStart=/usr/local/bin/qdgz300-spool-mover.sh
Restart=always
RestartSec=5
NoNewPrivileges=true
Nice=10
IOSchedulingClass=idle
CPUAffinity=0-15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
