# deploy/systemd/qdgz300-receiver.service
[Unit]
Description=QDGZ300 Receiver Service
After=network-online.target qdgz300-sysctl.service nic-optimization.service cpu-performance.service
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@INSTALL_ROOT@
Environment=LD_LIBRARY_PATH=/usr/local/corex/lib64:/usr/local/corex/lib
ExecStart=@INSTALL_ROOT@/bin/receiver_app --config @INSTALL_ROOT@/config/receiver.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
