# deploy/systemd/qdgz300-spool-mover.service
[Unit]
Description=QDGZ300 PCAP spool archive mover
After=network-online.target @RECEIVER_SERVICE@
Wants=network-online.target @RECEIVER_SERVICE@

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@INSTALL_ROOT@
Environment=QDGZ300_RECEIVER_CONFIG=@INSTALL_ROOT@/config/receiver.yaml
Environment=QDGZ300_DEFAULT_METRICS_DIR=@INSTALL_ROOT@/data/metrics
ExecStart=/usr/local/bin/qdgz300-spool-mover.sh
Restart=always
RestartSec=5
NoNewPrivileges=true
Nice=10
IOSchedulingClass=idle
CPUAffinity=0-15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
